# trading-bot/systemd/trading-bot.service
[Unit]
Description=Autonomous Trading Bot Service
After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=tradingbot
Group=tradingbot
WorkingDirectory=/mnt/c/Users/katbo/Desktop/Turbo Bot Deva/trading-bot
ExecStart=/usr/bin/npx ts-node autonomous_trading_bot.ts
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trading-bot

# Environment variables
Environment=NODE_ENV=production
Environment=BOT_MODE=demo
Environment=INITIAL_CAPITAL=10000
Environment=TRADING_INTERVAL=30000
Environment=KAFKA_ENABLED=true
Environment=KAFKA_BROKERS=localhost:9092
Environment=MAX_DRAWDOWN=0.15
Environment=MAX_DAILY_DRAWDOWN=0.05
Environment=AUTO_HEDGING=true
Environment=EXECUTION_ENGINE=simulated

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/mnt/c/Users/katbo/Desktop/Turbo Bot Deva/trading-bot

[Install]
WantedBy=multi-user.target
